# scripts/wasmstore.service
[Unit]
Description=wasmstore

[Service]
ExecStart=/usr/local/bin/wasmstore server
# WorkingDirectory=/
